# etc/systemd/system/systemd-networkd.service.d/wan-secret.conf
[Service]
LoadCredentialEncrypted=systemd.networkd.network.wan
