# etc/systemd/system/var.mount
[Mount]
What=LABEL=data
TimeoutSec=5min

[Install]
WantedBy=local-fs.target
